# queue-worker.service
[Unit]
Description=Trusted Advisors Queue Worker
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/trusted-advisors
ExecStart=/usr/bin/php artisan queue:work database --sleep=3 --tries=3 --max-time=3600
Restart=always
RestartSec=5
StandardOutput=append:/var/log/trusted-advisors-queue.log
StandardError=append:/var/log/trusted-advisors-queue.log

[Install]
WantedBy=multi-user.target
